# test/crontab
# Calculate 1024 + 1024 once per second
* * * * * * echo ABC${ARG}
