# test/crontab
# Construct string

* * * * * * echo ABC${ARG}
  
